# .config/autostart/Compositor.desktop
[Desktop Entry]
Exec=xcompmgr -c -C -t-5 -l-5 -r4.2 -o.55 &
Name=Compositor
Type=Application
Version=1.0
