# .config/autostart/Compositor.desktop
[Desktop Entry]
Exec=xcompmgr -n -f -c
Name=Compositor
Type=Application
Version=1.0
