# .config/autostart/xfce4-settings daemon.desktop
[Desktop Entry]
Exec=xfsettingsd &
Name=xfce4-settings daemon
Type=Application
Version=1.0
